# supervisor/conf.d/uwsgi.conf
[program:uwsgi]
command=/usr/local/pyenv/versions/2.7.14/bin/uwsgi --emperor "/etc/uwsgi/apps-enabled/*.ini" --die-on-term --master --uid www-data --gid www-data --logto /var/log/uwsgi/emperor.log

